# luci-app-uugamebooster/root/etc/uci-defaults/45_luci-uuplugin
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@uuplugin[-1]
	add ucitrack uuplugin
	set ucitrack.@uuplugin[-1].init=uuplugin
	commit ucitrack
EOF
chmod +x /etc/init.d/uuplugin
rm -f /tmp/luci-indexcache
exit 0
